# supervisor.d/icinga.conf
# Supervisor config for icinga daemon

[program:icinga]
command = %(here)s/conf.d/icinga.sh
stdout_events_enabled = true
stderr_events_enabled = true
