# supervisor.d/ido2db.conf
# Supervisor config for ido2db daemon

[program:ido2db]
command = %(here)s/conf.d/ido2db.sh
stdout_events_enabled = true
stderr_events_enabled = true
